# com.macmonica.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.macmonica.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/Caskroom/miniforge/base/bin/python3</string>
        <string>-m</string>
        <string>macmonica</string>
        <string>collect</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/hritikvalluvar/Projects/Tools/sysmon</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/hritikvalluvar/.macmonica/collector.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hritikvalluvar/.macmonica/collector-err.log</string>
</dict>
</plist>
